# deploy/ibtrading-backup.timer
[Unit]
Description=Daily ib-trading journal backup

[Timer]
OnCalendar=*-*-* 12:30:00
Persistent=true

[Install]
WantedBy=timers.target
