# deploy/ibtrading-flexrecon.timer
[Unit]
Description=Weekly Flex reconciliation (Saturday 20:30 server time)

[Timer]
OnCalendar=Sat 20:30
Persistent=true

[Install]
WantedBy=timers.target
