# deploy/ibtrading-weekendbrief.timer
[Unit]
Description=Weekend risk brief (Friday 17:00 America/New_York, DST-aware)

[Timer]
OnCalendar=Fri *-*-* 17:00:00 America/New_York
Persistent=true

[Install]
WantedBy=timers.target
